# files/etc/uci-defaults/99-luci-sso-ui
#!/bin/sh

# LuCI SSO - Semantic UI Integration
# Targeting LuCI 24.10+ (ucode templates)

# 1. CLEANUP: Remove any legacy patches from header.ut
CLEAN_HEADER() {
    FILE="$1"
    if [ -f "$FILE" ]; then
        sed -i '/luci-sso-login.js/d' "$FILE"
    fi
}

CLEAN_HEADER "/usr/share/ucode/luci/template/themes/bootstrap/header.ut"
CLEAN_HEADER "/usr/share/ucode/luci/template/header.ut"

# 2. SEMANTIC PATCH: Inject into sysauth.ut
PATCH_SYSAUTH() {
    FILE="$1"
    [ -f "$FILE" ] || return 0
    
    # Remove old patches
    sed -i '/luci-sso-login.js/d' "$FILE"
    
    # NEW STRATEGY: Use ABSOLUTE PATH to ensure loading
    sed -i "/{% include('header'.* %}/a \
<script src=\"/luci-static/resources/luci-sso-login.js\"></script>" "$FILE"
}

# Apply to global and bootstrap-specific templates
PATCH_SYSAUTH "/usr/share/ucode/luci/template/sysauth.ut"
PATCH_SYSAUTH "/usr/share/ucode/luci/template/themes/bootstrap/sysauth.ut"

# 3. Nuclear Cache Clear
rm -rf /tmp/luci-modulecache/* 2>/dev/null
rm -rf /tmp/luci-indexcache 2>/dev/null

exit 0